# salt/imgseek/upstart/iskdaemon.conf
description "imgseek"

start on runlevel [2345]
stop on runlevel [!2345]

script
exec {{ pillar['imgseek']['virtualenv_path'] }}/bin/iskdaemon.py
end script

respawn
