# salt/crontab/publisher
# every 4 hours
0 */4 * * *         {{ pillar['system']['user'] }}   {{ pillar['project']['virtualenv_path'] }}/bin/python {{ pillar['project']['path'] }}/manage.py post_feature_to_tumblr
